# systemd/update-beeper-user.timer
[Unit]
Description=Check for Beeper updates daily
Documentation=https://github.com/robertogogoni/update-beeper

[Timer]
# Run daily at a random time between 10:00 and 14:00
OnCalendar=*-*-* 10:00:00
RandomizedDelaySec=4h
# Run on boot if last scheduled run was missed
Persistent=true

[Install]
WantedBy=timers.target
