# systemd/update-beeper.service
[Unit]
Description=Update Beeper Desktop to latest version
Documentation=https://github.com/robertogogoni/update-beeper
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-beeper --notify

# IMPORTANT: Edit these values for your system:
# - Change 'User' to your username
# - Change '1000' in DBUS path to your UID (run 'id -u' to find it)
User=%u
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Note: For most users, the systemd USER service (update-beeper-user.service)
# is recommended instead of this system service. It automatically uses
# the correct user context without manual configuration.

[Install]
WantedBy=multi-user.target
